# node-demo-app/debian/postinst
#!/bin/sh
systemctl enable node-demo-app
systemctl start node-demo-app
